# scripts/nse.service
[Unit]
Description=Network Sandbox Engine Daemon
After=network.target

[Service]
Type=simple
User=root
# Starts the production server binding to the Unix socket /run/nse.sock
ExecStart=/usr/local/bin/nse serve
Restart=always
RestartSec=3

# Limit capabilities to network admin/raw socket injection for security hardening
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

[Install]
WantedBy=multi-user.target
